# examples/example.service
# Example systemd service unit file
# This file demonstrates various systemd directives and can be used for testing the LSP
[Unit]
Description=Example web server service
Documentation=https://example.com/docs
After=network.target
Wants=network-online.target
Before=multi-user.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/example
ExecStart=/usr/bin/example-server --port 8080 --config /etc/example/config.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=-/etc/default/example-server
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/example/logs
PrivateTmp=true

[Install]
WantedBy=multi-user.target
Also=example-server.socket
