# examples/example.service
# Example systemd service unit file
# This file demonstrates various systemd directives and can be used for testing the LSP
[Unit]
Description=Example web server service
Documentation=https://example.com/docs
After=network.target
Wants=network-online.target
Before=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/example-server --port 8080
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
Also=example-server.socket
